# bin/systemd/alsactl_restore.service
[Unit]
Description=Restore alsa settings

[Service]
ExecStart=/usr/bin/alsactl -f /home/fenton/.asound.state restore

[Install]
WantedBy=multi-user.target
